# files/dotfiles/systemd/user/spotify-notify.service
[Unit]
Description=Spotify Track Change Notifications
After=graphical-session.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'playerctl -p spotify -F metadata --format "{{artist}} - {{title}}" | while read line; do /home/craig/.local/bin/spotify-notify; done'
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
